# system/skeleton/etc/profile.d/0002gdk.sh
[ -d /run/xdg ] ||
	mkdir /run/xdg && chmod 700 /run/xdg

export XDG_RUNTIME_DIR=/run/xdg
